# deploy/systemd/hyperoil.service
[Unit]
Description=HyperOil v2 Pair Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hyperoil
Group=hyperoil
WorkingDirectory=/opt/hyperoil
ExecStart=/usr/bin/docker compose up --no-build
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/hyperoil/data
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
